# etc/systemd/system/burnin@.service
[Unit]
Description=Burn In HDD
BindsTo=dev-%i.device

[Service]
Type=simple
ExecStart=/bin/sh -c "exec /usr/local/bin/burnin %I"
